# portal-gun/setup/shares/service/bridge1.router/base-files/etc/local.d/ping-clients.start
#!/bin/sh

## TODO: Future iterations, like go should obviously
## only setup the one for the right verison of the router

# bridge1
nohup ping -c 1 10.1.1.2 &
nohup ping -c 1 10.1.1.10 &
nohup ping -c 1 10.1.1.11 &
nohup ping -c 1 10.1.1.12 &
# bridge11
nohup ping -c 1 10.1.1.2 &
nohup ping -c 1 10.1.1.10 &
nohup ping -c 1 10.1.1.11 &
nohup ping -c 1 10.1.1.12 &
# bridge12
nohup ping -c 1 10.1.1.2 &
nohup ping -c 1 10.1.1.10 &
nohup ping -c 1 10.1.1.11 &
nohup ping -c 1 10.1.1.12 &
