# portal-gun/machine-types/service-type/router-type/installation/tor.router.portal/base-files/etc/local.d/fix-dns.start
#!/bin/sh

echo "nameserver 1.1.1.1" >  /etc/resolv.conf
echo "nameserver 9.9.9.9" >> /etc/resolv.conf
echo "nameserver 8.8.8.8" >> /etc/resolv.conf
echo "nameserver 8.8.4.4" >> /etc/resolv.conf

sleep 4
