# portal-gun/machine-types/service-type/router-type/installation/tor.router.portal/base-files/etc/local.d/ping-clients.start
#!/bin/sh

###############################################################################
# TODO                                                                        #
# Future iterations, like go should obviously only setup the one for the      #
# right verison of the router.                                                #
###############################################################################

## Universe 0
#nohup ping -c 1 10.0.0.2 &
#nohup ping -c 1 10.0.0.10 &
#nohup ping -c 1 10.0.0.11 &
#nohup ping -c 1 10.0.0.12 &

#nohup ping -c 1 10.1.1.2 &
#nohup ping -c 1 10.1.1.10 &
#nohup ping -c 1 10.1.1.11 &
#nohup ping -c 1 10.1.1.12 &

## Universe 1
nohup ping -c 1 10.5.5.2 &
nohup ping -c 1 10.5.5.10 &
nohup ping -c 1 10.5.5.11 &
nohup ping -c 1 10.5.5.12 &

nohup ping -c 1 10.6.6.2 &
nohup ping -c 1 10.6.6.10 &
nohup ping -c 1 10.6.6.11 &
nohup ping -c 1 10.6.6.12 &

## Universe 2
#nohup ping -c 1 10.10.10.2 &
#nohup ping -c 1 10.10.10.10 &
#nohup ping -c 1 10.10.10.11 &
#nohup ping -c 1 10.10.10.12 &

#nohup ping -c 1 10.11.11.2 &
#nohup ping -c 1 10.11.11.10 &
#nohup ping -c 1 10.11.11.11 &
#nohup ping -c 1 10.11.11.12 &

# Universe 3
#nohup ping -c 1 10.15.15.2 &
#nohup ping -c 1 10.15.15.10 &
#nohup ping -c 1 10.15.15.11 &
#nohup ping -c 1 10.15.15.12 &
#
#nohup ping -c 1 10.16.16.2 &
#nohup ping -c 1 10.16.16.10 &
#nohup ping -c 1 10.16.16.11 &
#nohup ping -c 1 10.16.16.12 &
