# portal-gun/portals/service.portals/galaxy.router.portal/base-files/etc/local.d/fix-dns.start
#!/bin/sh

cat "nameserver 8.8.4.4" > /etc/resolv.conf
service openvpn restart
